# systemd/btcconncount.timer
[Unit]
Description=Hourly check of bitcoind connection count.

[Timer]
OnCalendar=hourly

[Install]
WantedBy=multi-user.target

